# systemd_files/burger_clear_sessions.service
[Unit]
Description=Clearing django's sessions


[Service]
WorkingDirectory=/opt/burger2/