# systemd_files/cert_renewal.service
[Unit]
Description=Cert Renewal

[Service]
ExecStart=sudo acme-nginx -d zed-chi.fun -d www.zed-chi.fun